# deploydeck.service
[Unit]
Description=DeployDeck - Your container deployment command center
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=docker
WorkingDirectory=/opt/deploydeck
ExecStart=/usr/local/bin/deploydeck --config /opt/deploydeck/config.yaml
Restart=always
RestartSec=5

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deploydeck

# Environment
Environment="DEPLOYDECK_LOG_LEVEL=info"

[Install]
WantedBy=multi-user.target
